# bot.conf
[program:sticktgbot]

directory=/home/ubuntu/sticktgbot
command=python3 app.py

autostart=true
autorestart=true
environment=HOME="/home/ubuntu",USER="ubuntu"

stderr_logfile=/home/ubuntu/sticktgbot/logfile.log